# supervisor/supervisord_sickbeard.conf
[program:sickbeard]
directory=/opt/sickbeard
command=python SickBeard.py
user = root
autostart=true
autorestart=true